# roles/common/files/node-exporter.service
[Unit]
Description=node_exporter
Wants=network-online.target
After=