# roles/common/files/node-exporter.service
[Unit]
Description=node_exporter
Wants=network-online.target
After=network-online.target  

[Service]
Type=simple
Execstart=/opt/node_exporter/